# salt/supervisor/service.conf.jinja2
[program:{{ program_name }}]
command={{ command }}
user={{ user }}
directory={{ directory }}

stderr_logfile={{ stdout }}
stderr_logfile={{ stderr }}
